# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[program:fuseki]
command=/opt/fuseki/fuseki-server --loc=/fuseki-data/databases/vehiculo --port=3030 /vehiculo
directory=/opt/fuseki
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
startsecs=5
priority=10

[program:api]
command=uvicorn src.main:app --host 0.0.0.0 --port %(ENV_API_PORT)s
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
startsecs=3
startretries=5
priority=20
